# infra/scripts/systemd/heliogram-frontend.service
[Unit]
Description=HelioGram Frontend (Vite Preview)
After=network.target heliogram-backend.service
Wants=heliogram-backend.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/helio-platform/apps/heliogram/frontend
Environment=NODE_ENV=production
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npm run preview -- --host 127.0.0.1 --port 5050 --strictPort
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
