# infra/scripts/systemd/heliogram-backend.service
[Unit]
Description=HelioGram Django Backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/helio-platform/apps/heliogram/backend
EnvironmentFile=/opt/helio-platform/apps/heliogram/.env
ExecStartPre=/opt/helio-platform/apps/heliogram/backend/venv/bin/python manage.py migrate --noinput
ExecStart=/opt/helio-platform/apps/heliogram/backend/venv/bin/gunicorn heliogram_core.wsgi:application --bind 127.0.0.1:8010 --workers 1 --timeout 120
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

